# ops/systemd/user/stellaratlas-horizon.service
[Unit]
Description=StellarAtlas owned Horizon API
Requires=stellaratlas-horizon-postgres.service
After=stellaratlas-horizon-postgres.service

[Service]
Type=simple
WorkingDirectory=/home/observe/stellarbeat-data/horizon
EnvironmentFile=/home/observe/.config/stellaratlas/full-history.env
Environment=NETWORK=pubnet
Environment=PORT=8000
Environment=ADMIN_PORT=8001
Environment=INGEST=true
Environment=DISABLE_TX_SUB=true
Environment=HISTORY_RETENTION_COUNT=0
Environment=CAPTIVE_CORE_HTTP_PORT=11626
Environment=CAPTIVE_CORE_PEER_PORT=11625
Environment=CAPTIVE_CORE_STORAGE_PATH=/home/observe/stellarbeat-data/horizon/captive-core/pubnet
Environment=CAPTIVE_CORE_LOG_PATH=/home/observe/stellarbeat-data/horizon/logs/captive-core.log
Environment=STELLAR_CORE_BINARY_PATH=/home/observe/stellarbeat-data/stellar-core/bin/stellar-core
Environment=LD_LIBRARY_PATH=/home/observe/stellarbeat-data/stellar-core/runtime/usr/lib/x86_64-linux-gnu
ExecStart=/home/observe/stellarbeat-data/horizon/bin/horizon serve --apply-migrations
Restart=on-failure
RestartSec=10
TimeoutStopSec=120
KillSignal=SIGTERM
LimitNOFILE=65536
MemoryMax=32G

[Install]
WantedBy=default.target

